# infra/systemd/quant-trend-backup.timer
[Unit]
Description=Nightly Quant Trend Lab PostgreSQL backup

[Timer]
OnCalendar=*-*-* 03:30:00 Asia/Seoul
Persistent=true
RandomizedDelaySec=300

[Install]
WantedBy=timers.target
